# recipes-ynov/warp7-ynov/warp7-ynov/thermo.service
[Unit]
Description=Sd Server service
[Service]
Type=simple
ExecStart=/usr/bin/Thermo
[Install]
WantedBy=multi-user.target